# repository-deepin-ufpr-brazil/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
  configure)
    systemctl enable repository-to-ufpr.service
;;

esac
exit 0
